# youtube2bilibili.service
[Unit]
Description=YouTube to Bilibili Auto Uploader
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/youtube_to_bilibili
ExecStart=/usr/bin/xvfb-run -a /usr/bin/python3 /root/youtube_to_bilibili/main.py
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
